# magneto-web/config/supervisord.conf
[program:apache2]
command=service apache2 start
[program:redis-server]
command=service redis-server start


